# docker/php/cron/crontab
* * * * * root echo "cron works" >> /var/log/cron-test.log
